# proxy.service
[Unit]
Description=Transparent TCP Proxy
After=network-online.target
ConditionPathExists=/opt/proxy.py

[Service]
DynamicUser=yes
WorkingDirectory=/opt
Environment="PYTHONPYCACHEPREFIX=/tmp"
ExecStart=python3 /opt/proxy.py
MemoryHigh=25%

[Install]
WantedBy=multi-user.target
