# proxy.service
[Unit]
Description=Transparent TCP Proxy
After=network-online.target
ConditionPathExists=/opt/proxy/run.sh

[Service]
DynamicUser=true
TemporaryFileSystem=/var/cache
CacheDirectory=proxy
WorkingDirectory=/var/cache/proxy
Environment="PYTHONPYCACHEPREFIX=/var/cache/proxy"
ExecStart=/opt/proxy/run.sh
MemoryLow=1%
MemoryHigh=50%
CPUSchedulingPolicy=rr
CPUSchedulingPriority=1

[Install]
WantedBy=multi-user.target
