# proxy.service
[Unit]
Description=Transparent TCP Proxy
After=network-online.target
ConditionPathExists=/opt/proxy/proxy

[Service]
DynamicUser=true
WorkingDirectory=/opt/proxy
ExecStart=/opt/proxy/proxy

[Install]
WantedBy=multi-user.target
